# etc/rc.local
#!/bin/sh

cputurbo.sh $(cat /usr/local/etc/cputurbo)
exit 0
