# etc/systemd/system/cpu-booster.service
[Unit]
Description=Auto CPU booster

[Service]
ExecStart=/usr/local/bin/cpu-booster.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
